# luci-app-systemupdate/files/etc/uci-defaults/99-systemupdate
#!/bin/sh

uci -q set systemupdate.systemupdate.enabled='1'
uci -q set systemupdate.systemupdate.repo='zjt37/zheng'
uci commit systemupdate

exit 0
