# custom-settings/root/etc/uci-defaults/99-custom-gateway
#!/bin/sh
uci batch <<-EOF
	set network.lan.gateway='192.168.2.2'
	set network.lan.dns='223.5.5.5'
	commit network
EOF
exit 0
